# amniotic.service
[Unit]
Description = Amniotic
[Service]
Type = simple
Restart = always
ExecStart = %h/.local/bin/amniotic
[Install]
WantedBy = default.target